# systemd/ada-back.service
# Ada back — persistent conversation loop (Bun microagent).
#
# Install:
#   cp systemd/ada-back.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now ada-back
#
# Fail-fast by design (docs/PLAN.md §9.3): exits if perception-voice or
# presence-voice are unreachable; Restart=on-failure brings it back once
# they are. Wants= orders it after them when they're enabled.

[Unit]
Description=Ada Back - conversation loop
Wants=perception-voice.service voice.service
After=perception-voice.service voice.service

[Service]
Type=simple
WorkingDirectory=/workspace/ada/back
ExecStart=%h/.bun/bin/bun ada-back.mjs
Restart=on-failure
RestartSec=3
Environment=PULSE_SERVER=unix:/run/user/%U/pulse/native
# Tool commands live in ~/.local/bin (govee, laptop, ...) — the systemd
# user manager's default PATH is just /usr/local/bin:/usr/bin, which made
# desk_light and the work-laptop activities fail with command-not-found.
Environment=PATH=%h/.local/bin:%h/.bun/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
